# bin/flaskr.service
[Unit]
Description=Gunicorn instance to serve flaskr
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/meredith-deploy-playground/flaskr
Environment="PATH=/home/ubuntu/meredith-deploy-playground/.venv/bin"
ExecStart=/home/ubuntu/meredith-deploy-playground/.venv/bin/gunicorn --workers 3 --bind 0.0.0.0:3000 --pid /run/gunicorn/gunicorn.pid --timeout 120 hello:app
PIDFile=/run/gunicorn/gunicorn.pid

[Install]
WantedBy=multi-user.target