# bin/flaskr.service
[Unit]
Description=Gunicorn instance to serve flaskr
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/meredith-deployment-playground/
Environment="PATH=/home/ubuntu/meredith-deployment-playground/.venv/bin"
ExecStart=/home/ubuntu/meredith-deployment-playground/.venv/bin/gunicorn --workers 3 --bind 0.0.0.0:3000 -m 007 flaskr:app
ExecReload=/bin/kill -s HUP $MAINPID

[Install]
WantedBy=multi-user.target